# deploy/paper-reader.timer
[Unit]
Description=Run paper-reader daily at 00:00

[Timer]
OnCalendar=*-*-* 00:00:00
Persistent=true
RandomizedDelaySec=600

[Install]
WantedBy=timers.target
